# openwrt/files/etc/rc.d/S41servald
#!/bin/sh /etc/rc.common

START=41

start() {
  # Start servald
  mkdir -p $SERVALINSTANCE_PATH
  servald start
}

stop() {
  servald stop
}
