# openwrt/package/luci/i18n/chinese/dist/etc/uci-defaults/luci-i18n-chinese
#!/bin/sh
uci batch <<-EOF
	set luci.languages.zh-cn=Chinese
        commit luci
EOF
	
